# agent/atlas-agent.service
[Unit]
Description=Atlas Command Center metrics agent
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=luka
Environment="PATH=/home/atlas/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/local/bin/atlas-agent
# optional: ATLAS_AGENT_TOKEN / ATLAS_AGENT_PORT (file may be absent)
EnvironmentFile=-/etc/atlas-agent.env
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
